# model-server.service
[Service]
Type=simple
ExecStart=/usr/bin/node /opt/model-server/index.js

[Install]
WantedBy=multi-user.target
